# templates/shifter_imagegw.service.erb
[Unit]
Description=Shifter Image Manager
After=munge.service

[Service]
Type=simple
# for now it will run as root, considering a closed volatile system, will leave chowns commented for later reference
# User=shifter
# Group=shifter
# PrivateTmp=true
Environment=PYTHONPATH=<%= @shifter_udiroot_dir %>/lib/python2.7/site-packages/
# ExecStartPre=/usr/bin/chown shifter:shifter /var/log/shifter_imagegw
ExecStart=+/usr/bin/gunicorn -b 0.0.0.0:5000 --backlog 2048 --access-logfile=<%= @imagegw_log_dir %>/access.log --log-file=<%= @imagegw_log_dir %>/error.log shifter_imagegw.api:app
RestartSec=2
Restart=on-failure

[Install]
WantedBy=multi-user.target